# scripts/launchd/com.fhir-place.daily-doc-sync.plist
<?xml version="1.0" encoding="UTF-8"?>
<!-- Daily docs sync at 06:00 ET. See com.fhir-place.daily-pm-triage.plist for install steps. -->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key><string>com.fhir-place.daily-doc-sync</string>
  <key>ProgramArguments</key><array><string>__HOME__/src/fhir-place/scripts/local/daily-doc-sync.sh</string></array>
  <key>StartCalendarInterval</key><dict>
    <key>Hour</key><integer>6</integer>
    <key>Minute</key><integer>0</integer>
  </dict>
  <key>WorkingDirectory</key><string>__HOME__/src/fhir-place</string>
  <key>StandardOutPath</key><string>__HOME__/src/fhir-place/logs/launchd-daily-doc-sync.out</string>
  <key>StandardErrorPath</key><string>__HOME__/src/fhir-place/logs/launchd-daily-doc-sync.err</string>
  <key>RunAtLoad</key><false/>
</dict>
</plist>
