# scripts/launchd/com.fhir-place.pr-fixup-dispatch.plist
<?xml version="1.0" encoding="UTF-8"?>
<!--
  PR fixup dispatcher at 09:00 and 14:00 ET (twice daily). Runs
  BEFORE the issue-mode engineer-dispatch in each window
  (engineer-dispatch is at 09:30 / 14:30) so red-CI / unresolved-thread
  / uat: needs-changes PRs are addressed before new work gets pulled
  from the backlog. The 30-minute stagger keeps both runs off the
  same Claude Max OAuth session simultaneously (rate-limit avoidance).
  See `docs/prompts/pr-fixup-dispatch.md` for what it does.

  See com.fhir-place.daily-pm-triage.plist for install steps.
-->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key><string>com.fhir-place.pr-fixup-dispatch</string>
  <key>ProgramArguments</key><array><string>__HOME__/src/fhir-place/scripts/local/pr-fixup-dispatch.sh</string></array>
  <key>StartCalendarInterval</key>
  <array>
    <dict>
      <key>Hour</key><integer>9</integer>
      <key>Minute</key><integer>0</integer>
    </dict>
    <dict>
      <key>Hour</key><integer>14</integer>
      <key>Minute</key><integer>0</integer>
    </dict>
  </array>
  <key>WorkingDirectory</key><string>__HOME__/src/fhir-place</string>
  <key>StandardOutPath</key><string>__HOME__/src/fhir-place/logs/launchd-pr-fixup-dispatch.out</string>
  <key>StandardErrorPath</key><string>__HOME__/src/fhir-place/logs/launchd-pr-fixup-dispatch.err</string>
  <key>RunAtLoad</key><false/>
</dict>
</plist>
