# scripts/launchd/com.fhir-place.hourly-engineer-dispatch.plist
<?xml version="1.0" encoding="UTF-8"?>
<!--
  Engineer dispatch at 09:00 and 14:00 ET (twice daily).

  Despite the file's "hourly" name (kept for path compatibility with the
  existing dispatch-engineer.sh shim and the docs/prompts/ filename),
  the plist fires twice a day while we get comfortable with the local
  drivers. Move back to hourly by replacing the StartCalendarInterval
  array below with a single `<dict><key>Minute</key><integer>30</integer></dict>`
  once Daniel signs off.

  See com.fhir-place.daily-pm-triage.plist for install steps.
-->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key><string>com.fhir-place.hourly-engineer-dispatch</string>
  <key>ProgramArguments</key><array><string>__HOME__/src/fhir-place/scripts/local/engineer-dispatch.sh</string></array>
  <key>StartCalendarInterval</key>
  <array>
    <dict>
      <key>Hour</key><integer>9</integer>
      <key>Minute</key><integer>0</integer>
    </dict>
    <dict>
      <key>Hour</key><integer>14</integer>
      <key>Minute</key><integer>0</integer>
    </dict>
  </array>
  <key>WorkingDirectory</key><string>__HOME__/src/fhir-place</string>
  <key>StandardOutPath</key><string>__HOME__/src/fhir-place/logs/launchd-engineer-dispatch.out</string>
  <key>StandardErrorPath</key><string>__HOME__/src/fhir-place/logs/launchd-engineer-dispatch.err</string>
  <key>RunAtLoad</key><false/>
</dict>
</plist>
