# scripts/launchd/com.fhir-place.pr-fixup-dispatch.plist
<?xml version="1.0" encoding="UTF-8"?>
<!--
  PR fixup dispatcher at 09:30 and 14:30 ET (twice daily).

  Staggered 30 minutes after the issue-mode engineer-dispatch (09:00 /
  14:00) so the issue-mode run finishes first — both share the engineer
  subagent and the same Claude Max OAuth session, so concurrent runs
  risk rate-limiting. See `docs/prompts/pr-fixup-dispatch.md` for what
  it does.

  See com.fhir-place.daily-pm-triage.plist for install steps.
-->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key><string>com.fhir-place.pr-fixup-dispatch</string>
  <key>ProgramArguments</key><array><string>__HOME__/src/fhir-place/scripts/local/pr-fixup-dispatch.sh</string></array>
  <key>StartCalendarInterval</key>
  <array>
    <dict>
      <key>Hour</key><integer>9</integer>
      <key>Minute</key><integer>30</integer>
    </dict>
    <dict>
      <key>Hour</key><integer>14</integer>
      <key>Minute</key><integer>30</integer>
    </dict>
  </array>
  <key>WorkingDirectory</key><string>__HOME__/src/fhir-place</string>
  <key>StandardOutPath</key><string>__HOME__/src/fhir-place/logs/launchd-pr-fixup-dispatch.out</string>
  <key>StandardErrorPath</key><string>__HOME__/src/fhir-place/logs/launchd-pr-fixup-dispatch.err</string>
  <key>RunAtLoad</key><false/>
</dict>
</plist>
