# scripts/launchd/com.fhir-place.hourly-uat-validation.plist
<?xml version="1.0" encoding="UTF-8"?>
<!-- Hourly UAT validation at :15. See com.fhir-place.daily-pm-triage.plist for install steps. -->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key><string>com.fhir-place.hourly-uat-validation</string>
  <key>ProgramArguments</key><array><string>__HOME__/src/fhir-place/scripts/local/hourly-uat-validation.sh</string></array>
  <key>StartCalendarInterval</key><dict><key>Minute</key><integer>15</integer></dict>
  <key>WorkingDirectory</key><string>__HOME__/src/fhir-place</string>
  <key>StandardOutPath</key><string>__HOME__/src/fhir-place/logs/launchd-uat-validation.out</string>
  <key>StandardErrorPath</key><string>__HOME__/src/fhir-place/logs/launchd-uat-validation.err</string>
  <key>RunAtLoad</key><false/>
</dict>
</plist>
